# ymir/backend/src/pymir-controller/supervisor/conf.d/controller.ini
[program:controller_server]
command=python controller/server.py -d -f controller/server_prd_config.yaml
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=./logs/controller.log        ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (0 means none, default 10)stdout_capture_maxbytes=1MB   ; number of bytes in 'capturemode' (default 0)
stdout_events_enabled=false   ; emit events on stdout writes (default false)
stdout_syslog=false           ; send stdout to syslog with process name (default false)
stderr_logfile=./logs/controller.log        ; stderr log path, NONE for none; default AUTO
stderr_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (0 means none, default 10)
stderr_capture_maxbytes=1MB   ; number of bytes in 'capturemode' (default 0)
stderr_events_enabled=false   ; emit events on stderr writes (default false)
stderr_syslog=false           ; send stderr to syslog with process name (default false)
stopasgroup=true
startsecs=50




[program:label_project_monitor]
command=python -u controller/label_project_monitor.py
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=./logs/label_project_monitor.log        ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (0 means none, default 10)stdout_capture_maxbytes=1MB   ; number of bytes in 'capturemode' (default 0)
stdout_events_enabled=false   ; emit events on stdout writes (default false)
stdout_syslog=false           ; send stdout to syslog with process name (default false)
stderr_logfile=./logs/label_project_monitor.log        ; stderr log path, NONE for none; default AUTO
stderr_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (0 means none, default 10)
stderr_capture_maxbytes=1MB   ; number of bytes in 'capturemode' (default 0)
stderr_events_enabled=false   ; emit events on stderr writes (default false)
stderr_syslog=false           ; send stderr to syslog with process name (default false)
stopasgroup=true
startsecs=50
