# ymir/auth/supervisor/conf.d/ymir_auth.ini
[program:app]
directory=/auth
command=/bin/bash -c "sh prestart.sh && sh start.sh"
numprocs=1
priority=99
autostart=true
autorestart=true
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/app_logs/ymir_auth.log        ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (0 means none, default 10)stdout_capture_maxbytes=1MB   ; number of bytes in 'capturemode' (default 0)
stdout_events_enabled=false   ; emit events on stdout writes (default false)
stdout_syslog=false           ; send stdout to syslog with process name (default false)
stderr_logfile=/app_logs/ymir_auth.log        ; stderr log path, NONE for none; default AUTO
stderr_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (0 means none, default 10)
stderr_capture_maxbytes=1MB   ; number of bytes in 'capturemode' (default 0)
stderr_events_enabled=false   ; emit events on stderr writes (default false)
stderr_syslog=false           ; send stderr to syslog with process name (default false)
stopasgroup=true
startsecs=50

